# ansible/roles/subscription-host/templates/vpn-bootstrap.service.j2
[Unit]
Description=VPN subscription delivery service (/sub/ + /bootstrap/)
Documentation=https://github.com/po4yka/vpn-deploy/blob/main/docs/SUBSCRIPTION-PLANE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vpn-bootstrap
Group=vpn-bootstrap
ExecStart=/usr/local/bin/vpn-bootstrap.py
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ subscription.subscription_dir | default('/var/lib/vpn-subscription') }} {{ subscription.reads_log | default('/var/log/vpn-subscription/reads.log') | dirname }}
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=yes
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
