# tests/snapshot/golden/watchdog/templates/vpn-watchdog.service.j2
[Unit]
Description=VPN watchdog probe
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/vpn-watchdog.env
ExecStart=/usr/local/sbin/vpn-watchdog.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
ReadWritePaths=/var/lib/vpn-watchdog
