# ansible/roles/watchdog/templates/vpn-watchdog.service.j2
[Unit]
Description=VPN watchdog probe
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/vpn-watchdog.env
ExecStart=/usr/local/sbin/vpn-watchdog.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths={{ watchdog.state_file | dirname }}
NoNewPrivileges=true
